# config/service/mariadb/run
#!/bin/bash

/usr/bin/mysqld_safe --skip-syslog > /dev/null 2>&1 & wait